# services/hermes-refactor/deploy/systemd/hermes-cron.service
# ===== systemd: hermes-cron.service + hermes-cron.timer — 采集定时任务 =====
# 部署方式同上

[Unit]
Description=Hermes 采集 cron（ScoreGG / TS / TK / 日报 / 心情）
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/life-os-frontend-v2/services/hermes-refactor/backend
EnvironmentFile=/home/ubuntu/life-os-frontend-v2/services/hermes-refactor/backend/.env

# 此处为 cron 入口——采集任务统一调度脚本
# 可扩展为跑单个采集：python3 -c "from app.collectors.scoregg import run; run()"
ExecStart=/home/ubuntu/life-os-frontend-v2/services/hermes-refactor/backend/venv/bin/python -c "from app.collectors.daily_report import run; run()"

StandardOutput=journal
StandardError=journal
