# services/hermes-refactor/deploy/systemd/hermes-jieyi.service
# ===== systemd: hermes-jieyi.service — 结衣后端 =====
# 部署路径：/etc/systemd/system/hermes-jieyi.service

[Unit]
Description=Hermes 结衣后端（知行合一数字分身 + 长期复盘）
Documentation=https://hermes-agent.nousresearch.com/docs
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/life-os-frontend-v2/services/hermes-refactor/backend

EnvironmentFile=/home/ubuntu/life-os-frontend-v2/services/hermes-refactor/backend/.env

ExecStart=/home/ubuntu/life-os-frontend-v2/services/hermes-refactor/backend/venv/bin/python -m uvicorn app.main:app \
  --host 127.0.0.1 \
  --port 8881 \
  --log-level info \
  --workers 1

Restart=always
RestartSec=10
LimitNOFILE=65536
MemoryMax=512M

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
